# app/playbooks/threat-groups/apt-g1041-sea-turtle.json
{
  "id": "apt-g1041",
  "num": 202,
  "name": "MITRE ATT&CK Group — Sea Turtle",
  "fullName": "Sea Turtle (G1041) Threat Group Hunt",
  "type": "Threat Group / APT Hunt",
  "severity": "Critical",
  "priority": "High",
  "detection": "Security Onion, Sysmon, OSQuery, Velociraptor, Elastic, Carbon Black",
  "scenario": "[Sea Turtle](https://attack.mitre.org/groups/G1041) is a Türkiye-linked threat actor active since at least 2017 performing espionage and service provider compromise operations against victims in Asia, Europe, and North America. [Sea Turtle](https://attack.mitre.org/groups/G1041) is notable for targeting registrars managing ccTLDs and complex DNS-based intrusions where the threat actor compromised DNS providers to hijack DNS resolution for ultimate victims, enabling [Sea Turtle](https://attack.mitre.org/groups/G1041) to spoof log in portals and other applications for credential collection. This default playbook uses MITRE ATT&CK mapped techniques and practical Security Onion, Sysmon, OSQuery, Velociraptor, Elastic, and Carbon Black hunts to investigate activity consistent with Sea Turtle.",
  "mitre": "T1078, T1078.003, T1133, T1190, T1199, T1566, T1059.004, T1203, T1505.003, T1557, T1074.002, T1114.001, T1213.006, T1560.001, T1071.001, T1685.006, T1690, T1583, T1583.001, T1583.002, T1583.003, T1584.002, T1588.002, T1588.004",
  "aliases": [
    "Sea Turtle",
    "Teal Kurma",
    "Marbled Dust",
    "Cosmic Wolf",
    "SILICON"
  ],
  "mitreGroupId": "G1041",
  "mitreUrl": "https://attack.mitre.org/groups/G1041",
  "tools": "Security Onion; Sysmon; OSQuery; Velociraptor; Elastic; Carbon Black",
  "sev": "critical",
  "cat": "Threat Groups",
  "source": "library",
  "updated": "2026-05-18",
  "detSteps": [
    {
      "n": 1,
      "title": "Profile Sea Turtle with MITRE ATT&CK context",
      "detail": "MITRE Group ID: G1041. Aliases: Sea Turtle, Teal Kurma, Marbled Dust, Cosmic Wolf, SILICON. Primary mapped tactics: Initial Access, Execution, Persistence, Privilege Escalation, Credential Access, Collection, Command and Control, Defense Impairment, Resource Development, Stealth. Mapped techniques: T1078 Valid Accounts, T1078.003 Local Accounts, T1133 External Remote Services, T1190 Exploit Public-Facing Application, T1199 Trusted Relationship, T1566 Phishing, T1059.004 Unix Shell, T1203 Exploitation for Client Execution, T1505.003 Web Shell, T1557 Adversary-in-the-Middle, T1074.002 Remote Data Staging, T1114.001 Local Email Collection, plus 15 additional mapped technique(s). Source: https://attack.mitre.org/groups/G1041. Use this step to scope the hunt, select relevant telemetry, and prioritize techniques that overlap the current alert or campaign.",
      "queries": {
        "security_onion": "# Security Onion KQL - Sea Turtle (G1041) ATT&CK technique pivots\n(rule.threat.technique.id:(T1078 OR T1078.003 OR T1133 OR T1190 OR T1199 OR T1566 OR T1059.004 OR T1203 OR T1505.003 OR T1557 OR T1074.002 OR T1114.001 OR T1213.006 OR T1560.001 OR T1071.001 OR T1685.006 OR T1690 OR T1583 OR T1583.001 OR T1583.002) OR threat.technique.id:(T1078 OR T1078.003 OR T1133 OR T1190 OR T1199 OR T1566 OR T1059.004 OR T1203 OR T1505.003 OR T1557 OR T1074.002 OR T1114.001 OR T1213.006 OR T1560.001 OR T1071.001 OR T1685.006 OR T1690 OR T1583 OR T1583.001 OR T1583.002) OR event.module:sysmon OR event.dataset:suricata.eve OR event.dataset:zeek.*)",
        "sysmon": "<Sysmon schemaversion=\"4.90\">\n  <EventFiltering>\n    <ProcessCreate onmatch=\"include\">\n      <Rule groupRelation=\"or\" name=\"Sea Turtle - ATT&CK process behavior\" technique=\"T1078,T1078.003,T1133,T1190,T1199,T1566,T1059.004,T1203\">\n        <CommandLine condition=\"contains any\">-enc;-encodedcommand;IEX;FromBase64String;DownloadString;regsvr32;mshta;wmic;winrs;rundll32</CommandLine>\n        <Image condition=\"end with any\">\\powershell.exe;\\cmd.exe;\\wscript.exe;\\cscript.exe;\\mshta.exe;\\regsvr32.exe;\\rundll32.exe;\\wmic.exe;\\winrs.exe</Image>\n      </Rule>\n    </ProcessCreate>\n    <ProcessAccess onmatch=\"include\">\n      <Rule groupRelation=\"and\" name=\"Sea Turtle - Credential access candidate\" technique=\"T1003\">\n        <TargetImage condition=\"is\">C:\\Windows\\System32\\lsass.exe</TargetImage>\n        <SourceImage condition=\"end with any\">\\procdump.exe;\\rundll32.exe;\\powershell.exe;\\mimikatz.exe</SourceImage>\n      </Rule>\n    </ProcessAccess>\n    <NetworkConnect onmatch=\"include\">\n      <Rule groupRelation=\"or\" name=\"Sea Turtle - C2 and lateral movement ports\" technique=\"T1071,T1021\">\n        <DestinationPort condition=\"is any\">53;80;443;445;3389;5985;5986;8080;8443</DestinationPort>\n      </Rule>\n    </NetworkConnect>\n    <RegistryEvent onmatch=\"include\">\n      <Rule groupRelation=\"or\" name=\"Sea Turtle - Persistence and defense evasion registry changes\" technique=\"T1547,T1562\">\n        <TargetObject condition=\"contains any\">\\Run;\\RunOnce;\\Services\\;\\Terminal Server\\WinStations\\RDP-Tcp\\UserAuthentication;\\Lsa\\Security Packages</TargetObject>\n      </Rule>\n    </RegistryEvent>\n    <DnsQuery onmatch=\"include\">\n      <QueryName condition=\"contains any\">.top;.xyz;.club;.online;.site;.cc</QueryName>\n    </DnsQuery>\n  </EventFiltering>\n</Sysmon>",
        "osquery": "SELECT pid, name, path, cmdline, parent, start_time\nFROM processes\nWHERE cmdline LIKE '%-enc%'\n   OR cmdline LIKE '%FromBase64String%'\n   OR cmdline LIKE '%DownloadString%'\n   OR name IN ('powershell.exe','cmd.exe','rundll32.exe','regsvr32.exe','mshta.exe','wmic.exe','winrs.exe');",
        "velociraptor": "LET procs = SELECT Name, Exe, CommandLine, Pid,\n  authenticode(filename=Exe).Trusted AS Trusted\nFROM pslist() WHERE Exe\nSELECT Name, Exe, CommandLine, Pid, Trusted\nFROM procs\nWHERE NOT Trusted OR CommandLine =~ '(?i)(-enc|frombase64string|downloadstring|regsvr32|mshta|wmic|winrs)'",
        "elastic": "# Elastic KQL - Sea Turtle (G1041) ATT&CK technique pivots\n(rule.threat.technique.id:(T1078 OR T1078.003 OR T1133 OR T1190 OR T1199 OR T1566 OR T1059.004 OR T1203 OR T1505.003 OR T1557 OR T1074.002 OR T1114.001 OR T1213.006 OR T1560.001 OR T1071.001 OR T1685.006 OR T1690 OR T1583 OR T1583.001 OR T1583.002) OR threat.technique.id:(T1078 OR T1078.003 OR T1133 OR T1190 OR T1199 OR T1566 OR T1059.004 OR T1203 OR T1505.003 OR T1557 OR T1074.002 OR T1114.001 OR T1213.006 OR T1560.001 OR T1071.001 OR T1685.006 OR T1690 OR T1583 OR T1583.001 OR T1583.002) OR event.module:sysmon OR event.dataset:suricata.eve OR event.dataset:zeek.*)",
        "carbon_black": "process_name:(powershell.exe OR cmd.exe OR rundll32.exe OR regsvr32.exe OR mshta.exe OR wmic.exe OR winrs.exe) OR netconn_count:[1 TO *]"
      }
    },
    {
      "n": 2,
      "title": "Hunt initial access, execution, and persistence behaviors",
      "detail": "Search for public-facing exploitation, suspicious script execution, LOLBins, web shells, scheduled tasks, service creation, and autorun registry changes. Tune by asset role: web servers, domain controllers, and admin jump boxes should have separate baselines.",
      "queries": {
        "security_onion": "# Initial access/execution/persistence hunt\n(event.category:process AND process.name:(powershell.exe OR cmd.exe OR wscript.exe OR cscript.exe OR mshta.exe OR regsvr32.exe OR rundll32.exe))\nOR (event.category:file AND file.extension:(aspx OR jsp OR php OR ps1 OR vbs OR js))\nOR (event.code:(4698 OR 4702 OR 7045))\nOR (registry.path:*\\Software\\Microsoft\\Windows\\CurrentVersion\\Run* OR registry.path:*\\SYSTEM\\CurrentControlSet\\Services\\*)",
        "sysmon": "<Sysmon schemaversion=\"4.90\">\n  <EventFiltering>\n    <ProcessCreate onmatch=\"include\">\n      <Rule groupRelation=\"or\" name=\"Sea Turtle - ATT&CK process behavior\" technique=\"T1078,T1078.003,T1133,T1190,T1199,T1566,T1059.004,T1203\">\n        <CommandLine condition=\"contains any\">-enc;-encodedcommand;IEX;FromBase64String;DownloadString;regsvr32;mshta;wmic;winrs;rundll32</CommandLine>\n        <Image condition=\"end with any\">\\powershell.exe;\\cmd.exe;\\wscript.exe;\\cscript.exe;\\mshta.exe;\\regsvr32.exe;\\rundll32.exe;\\wmic.exe;\\winrs.exe</Image>\n      </Rule>\n    </ProcessCreate>\n    <ProcessAccess onmatch=\"include\">\n      <Rule groupRelation=\"and\" name=\"Sea Turtle - Credential access candidate\" technique=\"T1003\">\n        <TargetImage condition=\"is\">C:\\Windows\\System32\\lsass.exe</TargetImage>\n        <SourceImage condition=\"end with any\">\\procdump.exe;\\rundll32.exe;\\powershell.exe;\\mimikatz.exe</SourceImage>\n      </Rule>\n    </ProcessAccess>\n    <NetworkConnect onmatch=\"include\">\n      <Rule groupRelation=\"or\" name=\"Sea Turtle - C2 and lateral movement ports\" technique=\"T1071,T1021\">\n        <DestinationPort condition=\"is any\">53;80;443;445;3389;5985;5986;8080;8443</DestinationPort>\n      </Rule>\n    </NetworkConnect>\n    <RegistryEvent onmatch=\"include\">\n      <Rule groupRelation=\"or\" name=\"Sea Turtle - Persistence and defense evasion registry changes\" technique=\"T1547,T1562\">\n        <TargetObject condition=\"contains any\">\\Run;\\RunOnce;\\Services\\;\\Terminal Server\\WinStations\\RDP-Tcp\\UserAuthentication;\\Lsa\\Security Packages</TargetObject>\n      </Rule>\n    </RegistryEvent>\n    <DnsQuery onmatch=\"include\">\n      <QueryName condition=\"contains any\">.top;.xyz;.club;.online;.site;.cc</QueryName>\n    </DnsQuery>\n  </EventFiltering>\n</Sysmon>",
        "osquery": "SELECT name, action, path FROM scheduled_tasks\nUNION\nSELECT name, path, status FROM services\nWHERE path LIKE '%AppData%' OR path LIKE '%Temp%' OR path LIKE '%ProgramData%';",
        "velociraptor": "SELECT FullPath, Mtime, Size\nFROM glob(globs=['C:/inetpub/wwwroot/**/*.aspx','C:/ProgramData/**/*.ps1','C:/Users/*/AppData/**/*.exe'])\nWHERE Mtime > timestamp(epoch=now().Unix - 7*24*60*60)",
        "elastic": "event.category:\"process\" and process.name:(\"powershell.exe\" or \"cmd.exe\" or \"mshta.exe\" or \"regsvr32.exe\" or \"rundll32.exe\")",
        "carbon_black": "(process_name:powershell.exe OR process_name:cmd.exe OR process_name:mshta.exe OR process_name:regsvr32.exe OR process_name:rundll32.exe) AND (cmdline:-enc OR cmdline:DownloadString OR cmdline:FromBase64String OR cmdline:http)"
      }
    },
    {
      "n": 3,
      "title": "Hunt credential access, discovery, and lateral movement",
      "detail": "Prioritize LSASS access, SAM/SECURITY hive reads, domain and network discovery, RDP/SMB/WinRM movement, and authentication anomalies. Correlate endpoint process telemetry with Zeek conn/smb/kerberos logs.",
      "queries": {
        "security_onion": "# Credential access + discovery + lateral movement\n(event.module:sysmon AND event.code:10 AND process.Ext.api.target_process.executable:*\\lsass.exe)\nOR (event.category:process AND process.command_line:(*nltest* OR *net group* OR *net view* OR *whoami /all* OR *ipconfig /all*))\nOR (event.category:network AND destination.port:(88 OR 135 OR 139 OR 389 OR 445 OR 3389 OR 5985 OR 5986))",
        "sysmon": "<Sysmon schemaversion=\"4.90\">\n  <EventFiltering>\n    <ProcessCreate onmatch=\"include\">\n      <Rule groupRelation=\"or\" name=\"Sea Turtle - ATT&CK process behavior\" technique=\"T1078,T1078.003,T1133,T1190,T1199,T1566,T1059.004,T1203\">\n        <CommandLine condition=\"contains any\">-enc;-encodedcommand;IEX;FromBase64String;DownloadString;regsvr32;mshta;wmic;winrs;rundll32</CommandLine>\n        <Image condition=\"end with any\">\\powershell.exe;\\cmd.exe;\\wscript.exe;\\cscript.exe;\\mshta.exe;\\regsvr32.exe;\\rundll32.exe;\\wmic.exe;\\winrs.exe</Image>\n      </Rule>\n    </ProcessCreate>\n    <ProcessAccess onmatch=\"include\">\n      <Rule groupRelation=\"and\" name=\"Sea Turtle - Credential access candidate\" technique=\"T1003\">\n        <TargetImage condition=\"is\">C:\\Windows\\System32\\lsass.exe</TargetImage>\n        <SourceImage condition=\"end with any\">\\procdump.exe;\\rundll32.exe;\\powershell.exe;\\mimikatz.exe</SourceImage>\n      </Rule>\n    </ProcessAccess>\n    <NetworkConnect onmatch=\"include\">\n      <Rule groupRelation=\"or\" name=\"Sea Turtle - C2 and lateral movement ports\" technique=\"T1071,T1021\">\n        <DestinationPort condition=\"is any\">53;80;443;445;3389;5985;5986;8080;8443</DestinationPort>\n      </Rule>\n    </NetworkConnect>\n    <RegistryEvent onmatch=\"include\">\n      <Rule groupRelation=\"or\" name=\"Sea Turtle - Persistence and defense evasion registry changes\" technique=\"T1547,T1562\">\n        <TargetObject condition=\"contains any\">\\Run;\\RunOnce;\\Services\\;\\Terminal Server\\WinStations\\RDP-Tcp\\UserAuthentication;\\Lsa\\Security Packages</TargetObject>\n      </Rule>\n    </RegistryEvent>\n    <DnsQuery onmatch=\"include\">\n      <QueryName condition=\"contains any\">.top;.xyz;.club;.online;.site;.cc</QueryName>\n    </DnsQuery>\n  </EventFiltering>\n</Sysmon>",
        "osquery": "SELECT * FROM process_open_sockets\nWHERE remote_port IN (88,135,139,389,445,3389,5985,5986);\n\nSELECT pid, name, cmdline FROM processes\nWHERE cmdline LIKE '%nltest%' OR cmdline LIKE '%net group%' OR cmdline LIKE '%whoami /all%' OR cmdline LIKE '%ipconfig /all%';",
        "velociraptor": "SELECT Pid, Process, LocalAddress, LocalPort, RemoteAddress, RemotePort, Status\nFROM netstat()\nWHERE RemotePort IN (88,135,139,389,445,3389,5985,5986)",
        "elastic": "event.category:\"network\" and destination.port:(88 or 135 or 139 or 389 or 445 or 3389 or 5985 or 5986)",
        "carbon_black": "netconn_port:(88 OR 135 OR 139 OR 389 OR 445 OR 3389 OR 5985 OR 5986) OR cmdline:(nltest OR \"net group\" OR \"whoami /all\")"
      }
    },
    {
      "n": 4,
      "title": "Hunt command-and-control, ingress transfer, and exfiltration",
      "detail": "Look for rare outbound destinations, DNS TXT/NULL or high-entropy queries, suspicious HTTP POST activity, proxy-like behavior, and large outbound transfers. Baseline by subnet and server role before suppressing.",
      "queries": {
        "security_onion": "# C2 and exfiltration hunt\n(event.category:network AND destination.port:(53 OR 80 OR 443 OR 8080 OR 8443 OR 1080 OR 8888))\nOR (event.dataset:zeek.dns AND (dns.question.type:(TXT OR NULL) OR dns.question.name:*[A-Za-z0-9]{20,}*))\nOR (event.dataset:zeek.http AND http.request.method:POST AND (network.bytes:>1000000 OR bytes:>1000000))\nOR (event.dataset:suricata.eve AND event.kind:alert)",
        "sysmon": "<Sysmon schemaversion=\"4.90\">\n  <EventFiltering>\n    <ProcessCreate onmatch=\"include\">\n      <Rule groupRelation=\"or\" name=\"Sea Turtle - ATT&CK process behavior\" technique=\"T1078,T1078.003,T1133,T1190,T1199,T1566,T1059.004,T1203\">\n        <CommandLine condition=\"contains any\">-enc;-encodedcommand;IEX;FromBase64String;DownloadString;regsvr32;mshta;wmic;winrs;rundll32</CommandLine>\n        <Image condition=\"end with any\">\\powershell.exe;\\cmd.exe;\\wscript.exe;\\cscript.exe;\\mshta.exe;\\regsvr32.exe;\\rundll32.exe;\\wmic.exe;\\winrs.exe</Image>\n      </Rule>\n    </ProcessCreate>\n    <ProcessAccess onmatch=\"include\">\n      <Rule groupRelation=\"and\" name=\"Sea Turtle - Credential access candidate\" technique=\"T1003\">\n        <TargetImage condition=\"is\">C:\\Windows\\System32\\lsass.exe</TargetImage>\n        <SourceImage condition=\"end with any\">\\procdump.exe;\\rundll32.exe;\\powershell.exe;\\mimikatz.exe</SourceImage>\n      </Rule>\n    </ProcessAccess>\n    <NetworkConnect onmatch=\"include\">\n      <Rule groupRelation=\"or\" name=\"Sea Turtle - C2 and lateral movement ports\" technique=\"T1071,T1021\">\n        <DestinationPort condition=\"is any\">53;80;443;445;3389;5985;5986;8080;8443</DestinationPort>\n      </Rule>\n    </NetworkConnect>\n    <RegistryEvent onmatch=\"include\">\n      <Rule groupRelation=\"or\" name=\"Sea Turtle - Persistence and defense evasion registry changes\" technique=\"T1547,T1562\">\n        <TargetObject condition=\"contains any\">\\Run;\\RunOnce;\\Services\\;\\Terminal Server\\WinStations\\RDP-Tcp\\UserAuthentication;\\Lsa\\Security Packages</TargetObject>\n      </Rule>\n    </RegistryEvent>\n    <DnsQuery onmatch=\"include\">\n      <QueryName condition=\"contains any\">.top;.xyz;.club;.online;.site;.cc</QueryName>\n    </DnsQuery>\n  </EventFiltering>\n</Sysmon>",
        "osquery": "SELECT p.pid, p.name, p.path, s.remote_address, s.remote_port\nFROM processes p\nJOIN process_open_sockets s ON p.pid = s.pid\nWHERE s.remote_port IN (53,80,443,8080,8443,1080,8888);",
        "velociraptor": "SELECT Pid, Process, RemoteAddress, RemotePort, Status\nFROM netstat()\nWHERE RemotePort IN (53,80,443,8080,8443,1080,8888)",
        "elastic": "event.category:\"network\" and destination.port:(53 or 80 or 443 or 8080 or 8443 or 1080 or 8888)",
        "carbon_black": "netconn_port:(53 OR 80 OR 443 OR 8080 OR 8443 OR 1080 OR 8888) AND NOT process_name:(chrome.exe OR msedge.exe OR firefox.exe OR outlook.exe OR teams.exe)"
      }
    }
  ],
  "contSteps": [
    {
      "title": "Contain affected hosts and identities",
      "detail": "Isolate confirmed compromised endpoints, disable impacted accounts, revoke active sessions/tokens, and block observed C2 destinations while preserving evidence.",
      "queries": {}
    },
    {
      "title": "Apply tactical network controls",
      "detail": "Block confirmed malicious infrastructure, restrict administrative protocols to jump boxes, and increase Security Onion alert visibility for the relevant MITRE techniques.",
      "queries": {}
    }
  ],
  "eradSteps": [
    {
      "title": "Remove persistence and actor tooling",
      "detail": "Remove malicious services, scheduled tasks, startup items, web shells, unauthorized accounts, and binaries identified during the hunt. Validate with Sysmon/OSQuery/Velociraptor before reconnecting hosts.",
      "queries": {}
    },
    {
      "title": "Patch exploited weaknesses",
      "detail": "Patch exploited public-facing applications, harden identity controls, rotate exposed credentials, and remediate vulnerable software or appliance firmware relevant to the observed technique set.",
      "queries": {}
    }
  ],
  "recSteps": [
    {
      "title": "Restore and monitor",
      "detail": "Restore systems from trusted backups where needed, re-enable network access in stages, and monitor Security Onion dashboards for recurrence of mapped techniques for at least two business cycles.",
      "queries": {}
    },
    {
      "title": "Improve ATT&CK coverage",
      "detail": "Update detection engineering backlog with uncovered ATT&CK techniques, tune noisy queries with environment-specific allowlists, and document lessons learned in the incident record.",
      "queries": {}
    }
  ]
}
